# deployment/raspberry/supervisor/notifier.conf
[program:smog-backend]
command=/home/pi/smog-backend/venv/bin/python notifier.py
directory=/home/pi/smog-backend
user=pi
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true